# systemd/cynagora-agent.socket.in
[Unit]
Requires=cynagora-check.socket
After=cynagora-check.socket

[Socket]
FileDescriptorName=agent
ListenStream=@DEFAULT_SOCKET_DIR@/cynagora.agent
SocketUser=@USER@
SocketGroup=@GROUP@
SocketMode=0660
SmackLabelIPIn=@
SmackLabelIPOut=@

Service=cynagora.service

[Install]
WantedBy=sockets.target
