# systemd/cynagora-agent.socket.in
[Socket]
FileDescriptorName=agent
ListenStream=@DEFAULT_SOCKET_DIR@/cynagora.agent
SocketMode=0600
SmackLabelIPIn=@
SmackLabelIPOut=@

Service=cynagora.service

[Unit]
Wants=cynagora.target
Before=cynagora.target

[Install]
WantedBy=sockets.target
