# systemd/cynagora-check.socket.in
[Unit]
Wants=cynagora.target
Before=cynagora.target

[Socket]
ExecStartPre=+/usr/bin/mkdir --context=System:Shared -p -m 755 @DEFAULT_SOCKET_DIR@
FileDescriptorName=check
ListenStream=@DEFAULT_SOCKET_DIR@/cynagora.check
SocketUser=@USER@
SocketGroup=@GROUP@
SocketMode=0666
SmackLabelIPIn=*
SmackLabelIPOut=@

Service=cynagora.service

[Install]
WantedBy=sockets.target
