# systemd/cynara-check.socket.in
[Socket]
FileDescriptorName=check
ListenStream=@SOCKET_DIR@/cynara.check
SocketMode=0666
SmackLabelIPIn=*
SmackLabelIPOut=@

Service=cynara.service

[Unit]
Wants=cynara.target
Before=cynara.target

[Install]
WantedBy=sockets.target
