# systemd/cynara-agent.socket.in
[Socket]
FileDescriptorName=agent
ListenStream=@DEFAULT_SOCKET_DIR@/cynara.agent
SocketMode=0600
SmackLabelIPIn=@
SmackLabelIPOut=@

Service=cynara.service

[Unit]
Wants=cynara.target
Before=cynara.target

[Install]
WantedBy=sockets.target
