# systemd/cynagora.target
[Unit]
Description=cynagora sockets
DefaultDependencies=true

