# systemd/cynagora.service.in
[Unit]
Description=Cynagora service

[Service]
ExecStart=/usr/bin/cynagorad --user @USER@ --group @GROUP@ --make-db-dir --own-db-dir

Type=notify

KillMode=process
TimeoutStopSec=3
Restart=always

Sockets=cynagora-admin.socket
Sockets=cynagora-check.socket
Sockets=cynagora-agent.socket

#NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
