# systemd/cynagora-check.socket.in
[Socket]
FileDescriptorName=check
ListenStream=@DEFAULT_SOCKET_DIR@/cynagora.check
SocketUser=@USER@
SocketGroup=@GROUP@
SocketMode=0666
SmackLabelIPIn=*
SmackLabelIPOut=@

Service=cynagora.service

[Unit]
Wants=cynagora.target
Before=cynagora.target

[Install]
WantedBy=sockets.target
